# packages/simo/simo-3.4.33.tar.gz/simo-3.4.33/src/simo/core/management/_hub_template/hub/supervisor.conf
# execute this on every program start
[program:simo-start]
directory={{ project_dir }}/hub/
command={{ venv_path }}/python manage.py on_http_start
process_name=%(program_name)s
user=root
stdout_logfile=/var/log/simo/start.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=3
redirect_stderr=true
autostart=true
autorestart=false
stopwaitsecs=15
killasgroup=true
environment=PYTHONUNBUFFERED=1


# using gunicorn for regular requests
[program:simo-gunicorn]
directory={{ project_dir }}/hub/
command={{ venv_path }}/gunicorn --workers 2 --timeout 120 --bind unix:/tmp/gunicorn.sock wsgi:application
process_name=%(program_name)s
user=root
stdout_logfile=/var/log/simo/gunicorn.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=3
redirect_stderr=true
autostart=true
autorestart=true
stopwaitsecs=15
killasgroup=true


# using daphne for socket connections routed to /ws/ on nginx.conf
[program:simo-daphne]
directory={{ project_dir }}/hub/
command={{ venv_path }}/daphne -u /tmp/http.sock --access-log /dev/stdout --proxy-headers asgi:application
process_name=%(program_name)s
user=root
stdout_logfile=/var/log/simo/daphne.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=3
redirect_stderr=true
autostart=true
autorestart=true
stopwaitsecs=15
killasgroup=true
environment=PYTHONUNBUFFERED=1


[program:simo-mcp]
directory={{ project_dir }}/hub/
command={{ venv_path }}/uvicorn simo.mcp_server.server:create_app --factory --host 0.0.0.0 --port 3333
process_name=%(program_name)s
user=root
stdout_logfile=/var/log/simo/mcp.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=3
redirect_stderr=true
autostart=true
autorestart=true
stopwaitsecs=15
killasgroup=true
environment=PYTHONUNBUFFERED=1


[program:simo-gateways]
command={{ venv_path }}/python {{ project_dir }}/hub/manage.py gateways_manager
process_name=%(program_name)s
user=root
stopwaitsecs=10
stdout_logfile=/var/log/simo/gateways.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=3
redirect_stderr=true
autostart=true
autorestart=true
stopwaitsecs=15
killasgroup=true
stopsignal=INT
environment=PYTHONUNBUFFERED=1


[program:simo-mqtt-control]
directory={{ project_dir }}/hub/
command={{ venv_path }}/python manage.py run_app_mqtt_control
process_name=%(program_name)s
user=root
stdout_logfile=/var/log/simo/mqtt_control.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=3
redirect_stderr=true
autostart=true
autorestart=true
stopwaitsecs=15
killasgroup=true
stopsignal=INT
environment=PYTHONUNBUFFERED=1

[program:simo-mqtt-fanout]
directory={{ project_dir }}/hub/
command={{ venv_path }}/python manage.py run_app_mqtt_fanout
process_name=%(program_name)s
user=root
stdout_logfile=/var/log/simo/mqtt_fanout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=3
redirect_stderr=true
autostart=true
autorestart=true
stopwaitsecs=15
killasgroup=true
stopsignal=INT
environment=PYTHONUNBUFFERED=1


[program:simo-celery-beat]
directory={{ project_dir }}/hub/
command={{ venv_path }}/celery -A celeryc.celery_app beat -l info --pidfile="/var/run/celerybeat.pid"
process_name=%(program_name)s
user=root
stdout_logfile=/var/log/simo/celery_beat.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=3
redirect_stderr=true
autostart=true
autorestart=true
stopwaitsecs=15
killasgroup=true


[program:simo-celery-worker]
directory={{ project_dir }}/hub/
command={{ venv_path }}/celery -A celeryc.celery_app worker -l info --concurrency=4
process_name=%(program_name)s
user=root
stdout_logfile=/var/log/simo/celery_worker.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=3
redirect_stderr=true
autostart=true
autorestart=true
stopwaitsecs=15
killasgroup=false
