# packages/eventix/eventix-4.3.3.tar.gz/eventix-4.3.3/supervisord.conf
[supervisord]
pidfile=/tmp/supervisord.pid
nodaemon=true
# user=nobody
# user=root
loglevel=INFO

[program:api]
priority=900
command=uvicorn main:app --host 0.0.0.0 --port %(ENV_EVENTIX_API_PORT)s --workers=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:worker]
priority=999
command=bash -c "sleep 10 && python worker.py"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
