# packages/sing-box-config/sing_box_config-0.8.0.tar.gz/sing_box_config-0.8.0/roles/sing_box_install/templates/etc/systemd/system/sing-box-reload.path.j2
[Unit]
Description=Watch {{ sing_box_etc_dir }}/config.json for changes

[Path]
Unit=sing-box-reload.service
PathChanged={{ sing_box_etc_dir }}/config.json

[Install]
WantedBy=multi-user.target
