# packages/hyodo/hyodo-3.0.1.tar.gz/hyodo-3.0.1/scripts/archive/com.afo.sync_phase3_issue_map.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>GH_PATH</key>
      <string>/opt/homebrew/bin/gh</string>
    </dict>
    <key>Label</key>
    <string>com.afo.sync_phase3_issue_map</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python3</string>
      <string>./scripts/sync_phase3_issue_map.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>WorkingDirectory</key>
    <string>.</string>
    <key>StandardOutPath</key>
    <string>./artifacts/run/sync_phase3_issue_map.out</string>
    <key>StandardErrorPath</key>
    <string>./artifacts/run/sync_phase3_issue_map.err</string>
  </dict>
</plist>
