# packages/bbperf/bbperf-0.0.32.tar.gz/bbperf-0.0.32/bbperf.service
[Unit]
Description=bbperf
After=syslog.target network.target

# StartLimitIntervalSec StartLimitBurst
#   stop restarting if fails X times with Y interval (zero disables)
#   default 5 times in 10 seconds
StartLimitIntervalSec=0

[Service]
Type=simple
User=mfreemon
ExecStart=/home/mfreemon/.venv-bbperf/bin/bbperf -s
Restart=always

# RestartSec
#   wait time until restart
RestartSec=2

[Install]
WantedBy=multi-user.target

