# packages/hyodo/hyodo-3.0.1.tar.gz/hyodo-3.0.1/scripts/cron/com.afo.meta_ssot_health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.afo.meta_ssot_health</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>./scripts/meta_ssot_health.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>.</string>

    <key>StartInterval</key>
    <integer>86400</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>./artifacts/run/meta_ssot_health.out</string>

    <key>StandardErrorPath</key>
    <string>./artifacts/run/meta_ssot_health.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
