# packages/data-safe-haven/data_safe_haven-5.7.0-py3-none-any.whl/data_safe_haven/resources/workspace/ansible/files/etc/systemd/system/clamav-clamonacc.service
[Unit]
Description=ClamAV on access scanning daemon
Requires=clamav-daemon.service
After=clamav-daemon.service syslog.target network.target

[Service]
Type=simple
User=root
ExecStart=/usr/sbin/clamonacc --foreground=true
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
