# packages/ratio1/ratio1-3.4.122.tar.gz/ratio1-3.4.122/ratio1/ipfs/ipfs_setup/ipfs.service
[Unit]
Description=IPFS daemon
After=network.target

[Service]
# Set the user who owns ~/.ipfs (assuming "ipfsuser" owns ~/.ipfs)
User=root
ExecStart=/usr/local/bin/ipfs daemon
Restart=always
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
