# packages/sing-box-config/sing_box_config-0.8.0.tar.gz/sing_box_config-0.8.0/roles/sing_box_config/templates/etc/systemd/system/sing-box-config-updater.timer.j2
[Unit]
Description=sing-box-config updater timer

[Timer]
OnBootSec=1min
OnUnitActiveSec={{ sing_box_config_updater_timer_interval }}
Persistent=true

[Install]
WantedBy=timers.target
